# var/spool/crontabs/radiocicletta
  0 4  *   *   *     /home/radiocicletta/aggiustone
